# server/devops/supervisor/supervisor.conf
[program:ws_server]
command = /www/subdomain.domain.com/poc-energy-meter-web/devops/bin/run_ws_server.sh
directory = /www/subdomain.domain.com/poc-energy-meter-web/server
user = www-data
autostart = true
autorestart = true
redirect_stderr = true
stdout_logfile = /var/log/supervisor/ws_server.log